# node-upstart.conf
description "node.js server"

# Start up when the system hits any normal runlevel, and
# shuts down when the system goes to shutdown or reboot.
start on startup
stop on shutdown

# IMPORTANT: You will want to use this with Forever. It
# tells Upstart that forever is going to fork after it
# starts.
expect fork

# This monitors Forever, which seems gratuitous.
# TIP: Comment these out while debugging your script.
respawn
respawn limit 5 60

# Send error messages to the console. Useful for debugging.
console output

# We run as
# this user instead of as root (the default user).
setuid ubuntu

# The user's home directory
env HOME=home/ubuntu

# Now we run Forever, telling it to write its logs to
# /tmp, and to run the script /opt/example/server.js
script
    cd $HOME
    exec forever start -a -l /tmp/forever.log -o /tmp/forever.stdout.log -e /tmp/forever.stderr.log /path/to/app.js
end script